# systemd/systemd/cvm-secrets-enhanced.service
[Unit]
Description=Enhanced cVM Secret Injection with Reboot Safety
After=network-online.target
Wants=network-online.target
Before=cvm-storage.service

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=VAULT_ADDR=https://vhsm.enclaive.cloud/
ExecStartPre=/usr/local/bin/vault-token-persistence.sh
ExecStart=/usr/local/bin/vhsm-cvm-auth-enhanced.sh

ExecStop=/bin/bash -c '\
    if [ -d /run/cvm-secrets ]; then \
        shred -vfz -n 3 /run/cvm-secrets/* 2>/dev/null || true; \
        rm -rf /run/cvm-secrets; \
    fi'

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target

